# root/etc/udev/rules.d/10-monitors.rules
ACTION=="add", SUBSYSTEM=="drm", RUN+="/bin/sh -c 'sleep 3 && /usr/bin/systemctl --user -M gavinr@ restart set-wallpaper.service'"
